# packaging/prerm
#!/bin/sh

systemctl disable grafsy.service
systemctl stop grafsy.service
